# configs/systemd/system/deploy-stacks.timer
[Unit]
Description=Deploy and update Docker stacks every hour 

[Timer]
OnCalendar=*-*-* *:00:00
Persistent=true

[Install]
WantedBy=timers.target